# api_db_testes.service
[Unit]
Description=dataCenterColect
[Service]
Type=simple
ExecStart=/etc/init.d/api_db_testes start

[Install]
WantedBy=multi-user.target
