# api_db_testes
#! /bin/sh
### BEGIN INIT INFO
# Provides:          api_db_testes.py
# Required-Start:    $syslog $time $remote_fs
# Required-Stop:     $syslog $time $remote_fs
# Default-Start:     3 4 5
# Default-Stop:      0 1 6
# Short-Description: api_db_testes.py
# Description:       api_db_testes
### END INIT INFO
#
# Author:      RCP
#


set -e

PATH=/bin:/usr/bin:/sbin:/usr/sbin:/usr/local/bin:/opt/api_db_testes
SL_PATH=/opt/api_db_testes
DAEMON="/opt/api_db_testes/api_db_testes.py "
PIDFILE=/var/run/api_db_testes.pid

test -x $SL_PATH || exit 0

. /lib/lsb/init-functions

case "$1" in
  start)
        log_daemon_msg "Starting deferred execution api_db_testes"
	start_daemon -p $PIDFILE $DAEMON
	#sudo $DAEMON
        log_end_msg $?
    ;;
  stop)
        log_daemon_msg "Stopping deferred execution api_db_testes"
        killproc -p $PIDFILE $DAEMON
        log_end_msg $?
    ;;
  force-reload|restart)
    $0 stop
    $0 start
    ;;
  status)
    status_of_proc -p $PIDFILE $DAEMON && exit 0 || exit $?
    ;;
  *)
    echo "Usage: /etc/init.d/api_db_testes {start|stop|restart|force-reload|status}"
    exit 1
    ;;
esac

exit 0
